# assets/flowy.service
[Unit]
Description=Flowy

[Service]
ExecStart=/usr/bin/flowy

[Install]
WantedBy=multi-user.target